# secret-santa.service
# ════════════════════════════════════════════════════════
#  Secret Santa Bot — systemd service fayli
#
#  O'rnatish tartibi:
#  1. sudo cp secret-santa.service /etc/systemd/system/
#  2. sudo systemctl daemon-reload
#  3. sudo systemctl enable secret-santa
#  4. sudo systemctl start secret-santa
#
#  Boshqaruv:
#    sudo systemctl status secret-santa    ← holat
#    sudo systemctl stop secret-santa      ← to'xtatish
#    sudo systemctl restart secret-santa   ← qayta ishga tushirish
#    journalctl -u secret-santa -f         ← real-time loglar
# ════════════════════════════════════════════════════════

[Unit]
Description=Secret Santa Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple

# --- BU QATORLARNI O'ZGARTIRASIZ ---
# Foydalanuvchi nomi (root dan foydalanmang!)
User=ubuntu

# Bot joylashgan papka to'liq yo'li:
WorkingDirectory=/home/ubuntu/secret_santa_bot

# Python virtual environment ishlatilsa:
ExecStart=/home/ubuntu/secret_santa_bot/venv/bin/python bot.py

# Virtual environment yo'q bo'lsa, bu qatorni ishlating:
# ExecStart=/usr/bin/python3 bot.py
# ────────────────────────────────────

# Xato bo'lganda avtomatik qayta ishga tushirish
Restart=always
RestartSec=10

# Muhit o'zgaruvchilari fayli
EnvironmentFile=/home/ubuntu/secret_santa_bot/.env

# Standart chiqish (journalctl da ko'rinadi)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secret-santa-bot

[Install]
WantedBy=multi-user.target
